# private_dot_config/systemd/user/backup_ha.timer
[Unit]
Description=Run Home Assistant backup sync daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
