# private_dot_config/systemd/user/backup_ha.service
[Unit]
Description=Backup Home Assistant backups from server

[Service]
Type=oneshot
ExecStart=%h/AutomatedScripts/back_ha.sh
