# private_dot_config/systemd/user/screendimming.timer
[Unit]
Description=Run gradual hyprsunset script every 5 minutes

[Timer]
# Run 30 seconds after boot, and then every 5 minutes
OnBootSec=30s
OnUnitActiveSec=5m
# Alternatively, use OnCalendar for exact minute alignment:
# OnCalendar=*:0/5
AccuracySec=1m # Allow systemd to slightly adjust timing for efficiency
Persistent=true # Run the timer job if the system was off during a scheduled time

[Install]
WantedBy=timers.target
