# private_dot_config/systemd/user/syncAndPublish-ObsidianToQuartz.timer
[Unit]
Description=Run sync-published-notes every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
